# deploy/acacia.service
[Unit]
Description=Acacia FastAPI Backend
After=network.target

[Service]
Type=simple
User=acacia
WorkingDirectory=/opt/acacia/backend
EnvironmentFile=/opt/acacia/backend/.env
ExecStart=/opt/acacia/backend/venv/bin/uvicorn main:app --host 127.0.0.1 --port 7860 --workers 2
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
